# config/systemd/user/pulseaudio.service.d/override.conf
[Service]
ExecStartPost=/bin/bash /home/simon/bin/move-mic-loopback.sh